# weclaw/service/weclaw.service
[Unit]
Description=WeClaw - WeChat AI Agent Bridge
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/weclaw start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
