# scripts/skillpay-auto-optimize.crontab
# SkillPay 收入优化自动化 - Crontab 配置示例
# 使用方法：crontab /Users/admin/.openclaw/workspace/scripts/skillpay-auto-optimize.crontab

# 环境变量
SHELL=/bin/zsh
PATH=/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin
HOME=/Users/admin

# ============================================
# 定时任务配置（取消注释以启用）
# ============================================

# 选项 1：每小时执行一次（推荐）
0 * * * * /Users/admin/.openclaw/workspace/scripts/skillpay-auto-optimize.sh >> /Users/admin/Library/Logs/skillpay-auto-optimize-cron.log 2>&1

# 选项 2：每 30 分钟执行一次
# */30 * * * * /Users/admin/.openclaw/workspace/scripts/skillpay-auto-optimize.sh >> /Users/admin/Library/Logs/skillpay-auto-optimize-cron.log 2>&1

# 选项 3：仅工作时间执行（9:00-22:00）
# 0 9-22 * * * /Users/admin/.openclaw/workspace/scripts/skillpay-auto-optimize.sh >> /Users/admin/Library/Logs/skillpay-auto-optimize-cron.log 2>&1

# 选项 4：仅工作日执行
# 0 * * * 1-5 /Users/admin/.openclaw/workspace/scripts/skillpay-auto-optimize.sh >> /Users/admin/Library/Logs/skillpay-auto-optimize-cron.log 2>&1

# 选项 5：每 2 小时执行一次
# 0 */2 * * * /Users/admin/.openclaw/workspace/scripts/skillpay-auto-optimize.sh >> /Users/admin/Library/Logs/skillpay-auto-optimize-cron.log 2>&1

# ============================================
# 管理命令
# ============================================
# 安装 crontab:  crontab /Users/admin/.openclaw/workspace/scripts/skillpay-auto-optimize.crontab
# 查看 crontab:  crontab -l
# 编辑 crontab:  crontab -e
# 删除 crontab:  crontab -r
# 查看 cron 日志: tail -f /Users/admin/Library/Logs/skillpay-auto-optimize-cron.log
